# 03-stage-wyoming-satellite/01-wyoming/files/wyoming-openwakeword.service
[Unit]
Description=Wyoming openWakeWord
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/wyoming-satellite/config.env
ExecStart=/opt/wyoming-openwakeword/bin/wyoming-openwakeword \
  --uri tcp://127.0.0.1:10400 \
  --preload-model ${WAKE_WORD}
Restart=on-failure
RestartSec=5
User=pi

[Install]
WantedBy=multi-user.target
