# 02-stage-audiodriver-2michat/01-driver/files/seeed-voicecard-install.service
[Unit]
Description=Install ReSpeaker 2-Mic HAT driver (seeed-voicecard) — first boot only
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/seeed-voicecard-installed

[Service]
Type=oneshot
ExecStart=/usr/local/bin/seeed-voicecard-install.sh
RemainAfterExit=yes
# Run at idle CPU/IO priority — keeps SSH and console responsive during DKMS build
Nice=19
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
